# missing-vowels.service
[Unit]
Description=Missing Vowels Game
After=network.target

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/missing-vowels/backend
Environment=PATH=/home/exedev/.local/bin:/usr/bin:/bin
ExecStart=/home/exedev/.local/bin/uv run python app.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
